# etc/systemd/user/imapsync-quick.service
[Unit]
Description=Offlineimap quick synchronization

[Service]
Type=simple
ExecStart=%h/.local/bin/imapsync -q
