# etc/systemd/user/imapsync-full.service
[Unit]
Description=Offlineimap full synchronization

[Service]
Type=simple
ExecStart=%h/.local/bin/imapsync
