# etc/systemd/user/mbsync.service
[Unit]
Description=Mail synchronization

[Service]
Type=oneshot
ExecStart=%h/.local/bin/mymailsync
ExecStartPost=/usr/bin/notmuch new
