# share/systemd/user/mailsync.service
[Unit]
Description=Run getmail

[Service]
Type=oneshot
ExecStart=/usr/bin/getmail --quiet

[Install]
WantedBy=default.target
