# share/systemd/user/mailsync.service
[Unit]
Description=Mail synchronization

[Service]
Type=oneshot
ExecStart=%h/.local/bin/mw
