# share/systemd/user/mgitdaemon.service
[Unit]
Description=mgitdaemon
After=network.target
[Service]
ExecStart=/home/roland/.local/bin/mgitdaemon
[Install]
WantedBy=default.target
