# etc/sh/profile.d/00-usrstuff.sh

# user data separate directories
# ~/.gnupg, ~/.password-store, ~/my and these ~/dotfiles
MY="${HOME}/my"
TASKDATA="${MY}/task"
FINANCE="${MY}/finance"

# environment variables usrstuff:
NAME="Roland Puntaier"
EMAIL="roland.puntaier@gmail.com"
MYGPG="F85CB631327ED8DCD50C9FD3BC9A552643B86228"

