# etc/sh/profile.d/00-usrstuff.sh

# user data separate directories
# ~/.gnupg, ~/.password-store, ~/my and these ~/dotfiles
export MY="${HOME}/my"
export TASKDATA="${MY}/task"
export FINANCE="${MY}/finance"

# environment variables usrstuff:
export NAME="Roland Puntaier"
export EMAIL="roland.puntaier@gmail.com"
export MYGPG="F85CB631327ED8DCD50C9FD3BC9A552643B86228"

