# etc/sh/profile.d/10-font.sh
FNTCONSOLE=/usr/share/kbd/consolefonts/ter-powerline-v14b.psf.gz
if [[ -f $FNTCONSOLE ]]; then
    setfont $FNTCONSOLE
fi
